# MacSettings/Library/LaunchAgents/com.local.Limactl.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>com.local.limactl</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/limactl</string>
        <string>hostagent</string>
        <string>--pidfile</string>
        <string>/Users/ronk/.lima/default/ha.pid</string>
        <string>--socket</string>
        <string>/Users/ronk/.lima/default/ha.sock</string>
        <string>default</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>UserName</key>
    <string>ronk</string>
    <key>WorkingDirectory</key>
    <string>/Users/ronk</string>
    <key>StandardOutPath</key>
    <string>/Users/ronk/Library/Logs/limactl/launch-limactl-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ronk/Library/Logs/limactl/launch-limactl-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin/:/usr/bin/</string>
    </dict>
</dict>
</plist>
